# profile.d/java.sh
#!/usr/bin/env bash

export JAVA_HOME=/usr/lib/jvm/default-jvm
export PATH=$PATH:/hbase/bin
export PATH=$PATH:$JAVA_HOME/bin
